# aws_workload_credentials_provider_common/configuration/aws-workload-credentials-provider-acm-nocaps.service
[Unit]
Description=AWS Workload Credentials Provider - ACM
After=network-online.target
Wants=network-online.target

[Service]
User=aws-workload-credentials-provider
WorkingDirectory=/opt/aws/workload-credentials-provider
Type=simple
Restart=on-failure
RestartSec=5s
TimeoutSec=1min
ExecStart=/opt/aws/workload-credentials-provider/bin/aws-workload-credentials-provider acm start

# Filesystem sandbox — /usr, /boot, /etc become read-only.
# Customer cert paths are added to ReadWriteDirectories via drop-in override
# generated by setup-config-based-permissions.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
LimitCORE=0

[Install]
WantedBy=multi-user.target
